# deploy/ansible/roles/proxy/templates/caddy.socket.j2
[Unit]
Description=Caddy socket listener

[Socket]
ListenStream=[::]:{{ baduk_port }}
BindIPv6Only=both

[Install]
WantedBy=sockets.target
RequiredBy=caddy.service
