# deploy/ansible/roles/db/templates/postgres-backup.timer.j2
[Unit]
Description=PostgreSQL Backup Timer
Requires=postgres.service

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
